# scheduled-tasks/com.schnapp.infra-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.schnapp.infra-health</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/scripts/check-infra-health.sh</string>
    </array>

    <!-- Run from the repo root (the probe reads no repo state, but keeps logs/relative paths sane). -->
    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <!-- launchd's default PATH lacks homebrew + docker/nc; prepend the dirs the probe needs. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Every 30 min (1800s) + once at load, so a downed service alerts within ~30 min, not a day. The
         probe is read-only and cheap; on a RED it opens an owner-assigned GitHub issue (email) via
         ops-alert.sh and exits non-zero (logged below). Alerting is transition-gated so a persistent RED
         is one issue, not spam. Substituted at install (see README): __REPO__ / __HOME__. -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/schnapp-os/infra-health.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/schnapp-os/infra-health.err.log</string>
</dict>
</plist>
