# scheduled-tasks/com.schnapp.session-mine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.schnapp.session-mine</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/scripts/session-mine-worker.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <!-- launchd's default PATH misses /usr/local/bin (claude) and homebrew git/gh/op. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <!-- op:// REFERENCE to the Claude credential (never the value); same headless-auth
             model as memory-consolidation (docs/headless-claude-auth.md, ADR 0019). -->
        <key>LEARNING_CLAUDE_TOKEN_REF</key>
        <string>__CLAUDE_TOKEN_REF__</string>
    </dict>

    <!-- Nightly at 03:40 local: after the 03:xx learning worker window, before the morning.
         Mining reads day-scale patterns; there is nothing to gain from firing more often. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>40</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/schnapp-os/session-mine.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/schnapp-os/session-mine.err.log</string>
</dict>
</plist>
