# scheduled-tasks/com.schnapp.caffeinate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.schnapp.caffeinate</string>

    <!-- Hold a power assertion so the hub Mac stays reachable (mac-mcp / obsidian-mcp / console
         tunnels) whenever it has power. `-s` = prevent system sleep, valid ONLY on AC power, so
         a MacBook on battery still sleeps normally (no battery-drain risk). This is a belt over
         `pmset -c sleep 0`: it survives a pmset change and covers clamshell-on-AC. It does NOT
         prevent sleep on battery, a manual Apple-menu Sleep, or the machine being powered off. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-s</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- caffeinate runs forever; KeepAlive restarts it if it is ever killed. -->
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/schnapp-os/caffeinate.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/schnapp-os/caffeinate.err.log</string>
</dict>
</plist>
