# scheduled-tasks/com.schnapp.memory-consolidation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.schnapp.memory-consolidation</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/plugins/core/scripts/learning-worker.sh</string>
    </array>

    <!-- Run from the repo root so the worker's own REPO_ROOT and relative paths resolve. -->
    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <!-- launchd's default PATH is only /usr/bin:/bin:/usr/sbin:/sbin, so `claude` (in
         /usr/local/bin) and homebrew git/gh would not be found and the live run would
         silently no-op. Prepend the dirs the worker + gate need. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <!-- op:// reference to the Claude OAuth token (a reference, not the secret). The worker
             resolves it via `op read` at runtime so the headless run authenticates despite
             launchd not reaching the login Keychain. Substituted at install (see README). -->
        <key>LEARNING_CLAUDE_TOKEN_REF</key>
        <string>__CLAUDE_TOKEN_REF__</string>
    </dict>

    <!-- Run nightly at 03:17 local time. RunAtLoad false: never runs automatically on launchctl load. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>17</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/schnapp-os/memory-consolidation.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/schnapp-os/memory-consolidation.err.log</string>
</dict>
</plist>
